# avahi-reflector/rootfs/etc/services.d/avahi-daemon/run
#!/usr/bin/env bash
# ==============================================================================
# Start avahi-daemon service
# ==============================================================================
exec avahi-daemon
